# hermit.service
[Unit]
Description=Hermit AI Agent OS
After=network.target docker.service
Requires=docker.service

[Service]
Type=simple
User=ubuntu
WorkingDirectory=/root/HermitShell
ExecStart=/root/HermitShell/hermit-server
Restart=always
RestartSec=10
StandardOutput=append:/root/HermitShell/data/logs/hermit.log
StandardError=append:/root/HermitShell/data/logs/hermit-error.log

[Install]
WantedBy=multi-user.target
